# master/supervisord.conf
[supervisord]
nodaemon=true

[program:saltmaster]
command=salt-master -l debug

[program:saltkey]
command= sh accept.sh